# templates/mongodb-mms-monitoring-agent.service.tmpl
[Unit]
Description=MongoDB monitoring agent

[Service]
Type=simple
User=root

ExecStart={{directory}}/mongodb-mms-monitoring-agent-script

[Install]
WantedBy=multi-user.target
